# configs/chimeratk-watchdog.service.in
[Unit]
Description=ChimeraTK Watchdog Server
After=network.target local-fs.target

[Service]
Type=simple
WorkingDirectory=/var/lib/chimeratk/@PROJECT_NAME@/
ExecStart=/usr/bin/watchdog-server
Restart=always
RestartSec=1min

[Install]
WantedBy=multi-user.target
