# configs/ChimeraTK-Watchdog.service.in
[Unit]
Description=ChimeraTK Watchdog Server
After=network.target local-fs.target

[Service]
Type=forking
ExecStart=@CMAKE_INSTALL_PREFIX/@PROJECT_NAME@/ChimeraTK-WatchdogServer
PIDFile=@CMAKE_INSTALL_PREFIX/@PROJECT_NAME@/WatchdogServer.PID
Restart=always
RestartSec=1min


[Install]
WantedBy=multi-user.target
